# openwrt/WR84X_qaul.net/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

## qaul
# autoconfigure
/qaul/firststartup.sh


exit 0
